# src/dist/ops/systemd/ccclaw-run.timer
[Unit]
Description=Run ccclaw worker on schedule (user)

[Timer]
OnCalendar=*:0/10 Asia/Shanghai
Persistent=true
Unit=ccclaw-run.service

[Install]
WantedBy=timers.target
