# src/dist/ops/systemd/ccclaw-ingest.timer
[Unit]
Description=Run ccclaw ingest on schedule (user)

[Timer]
OnCalendar=*:0/5 Asia/Shanghai
Persistent=true
Unit=ccclaw-ingest.service

[Install]
WantedBy=timers.target
