# src/dist/ops/systemd/ccclaw-run.service
[Unit]
Description=ccclaw run service
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ccclaw
ExecStart=/usr/local/bin/ccclaw run --config /opt/ccclaw/ops/config/config.toml --env-file /opt/ccclaw/.env
