# systemd/ccclaw-run.timer
[Unit]
Description=ccclaw run timer — 每 10 分钟执行一次任务

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s

[Install]
WantedBy=timers.target
