# src/ops/systemd/ccclaw-patrol.timer
[Unit]
Description=Run ccclaw patrol health/compensation on schedule (user)

[Timer]
OnCalendar=*:0/2 Asia/Shanghai
Persistent=true
Unit=ccclaw-patrol.service

[Install]
WantedBy=timers.target
