# src/ops/systemd/ccclaw-journal.timer
[Unit]
Description=Run ccclaw journal on schedule (user)

[Timer]
OnCalendar=*-*-* 23:50:00 Asia/Shanghai
Persistent=true
Unit=ccclaw-journal.service

[Install]
WantedBy=timers.target
