# src/ops/systemd/ccclaw-run.service
[Unit]
Description=ccclaw run service (user)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/.ccclaw
ExecStart=%h/.ccclaw/bin/ccclaw run --config %h/.ccclaw/ops/config/config.toml --env-file %h/.ccclaw/.env
