# src/ops/systemd/ccclaw-ingest.service
[Unit]
Description=ccclaw ingest service
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ccclaw
ExecStart=/usr/local/bin/ccclaw ingest --config /opt/ccclaw/ops/config/config.toml --env-file /opt/ccclaw/.env
